# 02-中型车个性化配置/配置服务/can.service
[Unit]
Description=can service for agilex
After=network.target

[Service]
Type=simple
ExecStart=/bin/sh -c "ip link set can0 up type can bitrate 500000"
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
